# config/etc/profile.d/docker.sh
#!/bin/bash

export DOCKER_TLS_VERIFY=1
export DOCKER_HOST=localhost:2376
export DOCKER_CERT_PATH=/data/docker/certs/$USER

